# deploy/config/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:gunicorn]
command=/opt/mysite/venv/bin/gunicorn config.wsgi:application -c /opt/mysite/deploy/config/gunicorn.conf.py
directory=/opt/mysite/backend
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/mysite/gunicorn.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=DJANGO_SETTINGS_MODULE="config.settings"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/mysite/nginx.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[eventlistener:process_exit]
command=/bin/bash -c "echo PROCESS_STATE_EXITED"
events=PROCESS_STATE_EXITED
